# scripts/launchd.nakano-monitor.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>com.nakano-yusaku.bot.monitor</string>

    <key>ProgramArguments</key>
    <array>
      <string>/bin/zsh</string>
      <string>-lc</string>
      <string>cd /Users/shonomasaki/Documents/Codex/codex-hermes-agent-x-search-gpt &amp;&amp; ./scripts/run-nakano-monitor.sh</string>
    </array>

    <key>StartInterval</key>
    <integer>3600</integer>

    <key>RunAtLoad</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/shonomasaki/Documents/Codex/codex-hermes-agent-x-search-gpt/outputs/launchd.nakano-monitor.out.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/shonomasaki/Documents/Codex/codex-hermes-agent-x-search-gpt/outputs/launchd.nakano-monitor.err.log</string>
  </dict>
</plist>
